# systemd/awg-bot.service
[Unit]
Description=AmneziaWG Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/amneziawg-bot
EnvironmentFile=/opt/amneziawg-bot/.env
ExecStart=/usr/bin/python3 /opt/amneziawg-bot/app/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target